# common/.config/systemd/user/llama-swap.service
[Unit]
Description=llama-swap AI inference router (GPU)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pannet1/programs/shell/github.com/pannet1/Scripts/ai
ExecStart=/usr/local/bin/llama-swap --config /home/pannet1/programs/shell/github.com/pannet1/Scripts/ai/config/llama-swap/config-host.yaml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
